# capture-image.service
[Unit]
Description=Capture image

[Service]
ExecStart=/usr/src/app/captureImage.sh

[Install]
WantedBy=multi-user.target


