# capture-image.service
[Unit]
Description=Capture image

[Service]
ExecStart=/bin/sh -c '/usr/src/app/captureImage.sh'

