# capture-image.service
[Unit]
Description=Capture image

[Service]
ExecStart=/bin/date >> /usr/src/app/sam
