# capture-image.service
[Unit]
Description=Capture image

[Service]
ExecStart=/bin/sh -c '/bin/date >> /usr/src/app/sam'
