# capture-image.service
[Unit]
Description=Capture image

[Service]
ExecStart=/bin/bash /usr/src/app/captureImage.sh

