# process-images.service
[Unit]
Description=Process images

[Service]
ExecStart=/bin/bash /usr/src/app/processImages.sh
