# process-images.service
[Unit]
Description=Process images

[Service]
ExecStart=/bin/sh -c '/usr/src/app/processImages.sh'
