# process-images.timer
[Unit]
Description=Process images every hour

[Timer]
OnBootSec=5min
OnCalendar=*-*-* *:00:00
