# process-images.service
[Unit]
Description=Process images

[Service]
ExecStart=/usr/src/app/processImages.sh

[Install]
WantedBy=multi-user.target

